# inProgress/1.Testing/Systemd/rclone-vfs.service
[Unit]
Description=rclone-vfs Container
Requires=docker.service
After=docker.service

[Service]
Restart=always
ExecStartPre=-/bin/fusermount -uz VFSPATH/rclone-vfs
ExecStart=/usr/bin/docker start -a rclone-vfs
ExecStop=/usr/bin/docker stop -t 2 rclone-vfs

[Install]
WantedBy=local.target
